# hms-firetv.service
[Unit]
Description=HMS FireTV - Fire TV Lightning Protocol Service
Documentation=https://github.com/yourusername/hms-firetv
After=network-online.target postgresql.service docker.service
Wants=network-online.target
Requires=postgresql.service

[Service]
Type=simple
User=aamat
Group=aamat
WorkingDirectory=/home/aamat/maestro_hub/hms-firetv/build

# Environment variables
Environment="DB_HOST=192.168.2.15"
Environment="DB_PORT=5432"
Environment="DB_NAME=firetv"
Environment="DB_USER=firetv_user"
Environment="DB_PASSWORD=firetv_postgres_2026_secure"
Environment="MQTT_BROKER_HOST=192.168.2.15"
Environment="MQTT_BROKER_PORT=1883"
Environment="MQTT_USER=aamat"
Environment="MQTT_PASS=exploracion"
Environment="API_HOST=0.0.0.0"
Environment="API_PORT=8888"
Environment="THREAD_NUM=4"
Environment="LOG_LEVEL=info"

# Service binary
ExecStart=/home/aamat/maestro_hub/hms-firetv/build/hms_firetv

# Restart policy
Restart=on-failure
RestartSec=5s
StartLimitInterval=60s
StartLimitBurst=3

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hms-firetv

[Install]
WantedBy=multi-user.target
